# systemd/mlgpt.service.in
[Unit]
Description=MLGPT (Streamlit)
Documentation=file://{{REPO}}/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Optional: create {{REPO}}/.env with e.g. MLGPT_CWD=/path/to/your/repo
EnvironmentFile=-{{REPO}}/.env
ExecStart=/usr/bin/bash {{REPO}}/scripts/run-streamlit-service.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
